# wiregate-INTERNAL/testing/elab17-wg/xPLLib/xplhub-1.31/DEBIAN/postinst
#! /bin/sh -e

update-rc.d xplhub defaults
/etc/init.d/xplhub restart
if [ -e /etc/init.d/monit ]
then
	/etc/init.d/monit restart
fi
