# wiregate-INTERNAL/testing/elab17-wg/irserver-6.02.23/debian/irserver/DEBIAN/postinst
#! /bin/sh -e

update-rc.d irserver defaults
/etc/init.d/irserver restart
